# yolo_fps_web.service
[Unit]
Description=零延迟YOLO FPS云辅助系统Web服务
After=network.target

[Service]
Type=simple
User=codespace
WorkingDirectory=/workspaces/zero-latency-yolo/src/web
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5
StandardOutput=append:/workspaces/zero-latency-yolo/logs/web_server.log
StandardError=append:/workspaces/zero-latency-yolo/logs/web_error.log
Environment="PORT=3000"

[Install]
WantedBy=multi-user.target
